# sys-config/user-config/files/etc/cron.monthly/defrag.monthly
#!/bin/sh
# Southern Tools
#
e4defrag -v /
e4defrag -v /home
