# sys-config/user-config/files/etc/cron.monthly/rsnapshot.monthly
#!/bin/sh
# Southern Tools
#
echo "### RSNAPSHOT MONTHLY ###"
mount -L BACKUP /mnt/backup && rsnapshot -V monthly || echo "Backup failure"
umount /mnt/backup
echo ""
