# sys-config/user-config/files/etc/cron.daily/genup
#!/bin/bash
# Southern Tools
#
export PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/bin"
genup -p >/var/log/latest-genup-run.log 2>&1
