# sys-config/southern-config/files_skel/etc/cron.monthly/rsnapshot.monthly
#!/bin/sh
# Southern Tools
#
echo "### RSNAPSHOT MONTHLY ###"
cryptsetup --key-file /root/.user_config/no_share/luks/encrypted_drive.key luksOpen /dev/disk/by-uuid/ENCRYPTED_DRIVE_UUID encrypted_drive
vgchange --activate y vg_encrypted
mount -L backup /mnt/backup && rsnapshot monthly || echo "Backup failure"
umount /mnt/backup
vgchange --activate n vg_encrypted
cryptsetup luksClose encrypted_drive
echo ""
