# sys-config/user-config/files_carson/etc/local.d/postboot.start
#!/bin/bash
# Southern Tools
#
# commands in this script will be run once the main boot process is complete
# ensure the framebuffer is properly cleared...
cat /dev/zero > /dev/fb0 2>/dev/null
