# sys-config/southern-config/files_skel/etc/local.d/umount_external_drive.stop
#!/bin/bash
# Southern Tools
#
umount /mnt/encrypted_unit_1
umount /mnt/encrypted_unit_2
umount /mnt/encrypted_unit_3
vgchange --available n vg_encrypted
cryptsetup luksClose encrypted_drive
