# sys-config/user-config/files_rand/etc/cron.daily/rsnapshot.daily
#!/bin/sh
# Southern Tools
#
echo "### RSNAPSHOT DAILY ###"
mount -L BACKUP /mnt/backup && rsnapshot -V daily || echo "Backup failure"
umount /mnt/backup
echo ""
