# sys-config/southern-config/files_skel/etc/cron.weekly/rsnapshot.weekly
#!/bin/sh
# Southern Tools
#
echo "### RSNAPSHOT WEEKLY ###"
mount -L BACKUP /mnt/backup && rsnapshot -V weekly || echo "Backup failure"
umount /mnt/backup
echo ""
