# sys-config/southern-config/files_skel/etc/cron.weekly/rsnapshot.weekly
#!/bin/sh
# Southern Tools
#
echo "### RSNAPSHOT WEEKLY ###"
cryptsetup --key-file /root/.user_config/no_share/luks/encrypted_drive.key luksOpen /dev/disk/by-uuid/ENCRYPTED_DRIVE_UUID encrypted_drive
vgchange --available y vg_encrypted
mount -L backup /mnt/backup && rsnapshot weekly || echo "Backup failure"
umount /mnt/backup
vgchange --available n vg_encrypted
cryptsetup luksClose encrypted_drive
echo ""
